# Applications/Crontab/cron_dir/a.crontab
#定时任务测试a
06 12 * * * root echo 'a' >> test.txt
* * * * * root echo 'a' >> test.txt
* * * * * root echo 'aca' >> test.txt