# Applications/Crontab/cron_dir/auto_Browser.crontab
#自动断点任务
* * * * * www  php Browser.php -u http%3A%2F%2Fwww.modulesoap.com%2FJump.php