# Applications/Crontab/cron_dir/b.crontab
#定时任务测试b
06 12 * * * root echo 'b' >> test.txt
